# units/.config/systemd/user/reflector.service
[Unit]
Description=Pacman mirrorlist update
